# services/offlineimap/run
#!/bin/execlineb

fdmove -c 2 1
offlineimap
