# services/log/run
#!/bin/dash

exec sudo -u user cyclog /opt/log/user <&0
