# services/cerberus/run
#!/bin/execlineb

fdmove -c 2 1
cerberus
