# services/tail/run
#!/bin/execlineb

tail -F /opt/log/offlineimap/current /opt/log/cerberus/current
